# vaultwarden/backup.cron
30 2 * * * root /lab/vaultwarden/backup.sh backup 2>&1 >>/lab/vaultwarden/backup/log
